# fs/skeleton/etc/inittab
console::respawn:-/etc/setupshell
console::sysinit:/bin/mount -n -a
console::once:/etc/rc.start
