# fs/skeleton/etc/inittab
null::sysinit:/bin/mount -t proc proc /proc
null::sysinit:/bin/mount -t jffs2 -o ro /dev/mtdblock2 /.flash

console::respawn:-/etc/setupshell
console::sysinit:/bin/mount -n -a
console::once:/init.sh
